# .config/autostart/RiseupVPN.desktop
[Desktop Entry]
Type=Application
Name=RiseupVPN
Exec="riseup-vpn"
X-GNOME-Autostart-enabled=true
